# tests/data/config1.ini
[supervisord]
tags = tag1:a, foo

[inet_http_server]
port = :19001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:p1]
command = cat
tags = tag2:b

[program:p2]
command = cat
tags = tag2:b, bar
numprocs = 3
process_name = p2_%(process_num)s

[program:p3]
command = cat
tags = tag2:c

[eventlistener:sup_broadcast]
command = sup_broadcast.py -c tests/data/config1.ini
events = PROCESS_STATE, SUPERVISOR_STATE_CHANGE, TICK_60, PROCESS_GROUP
buffer_size = 8192

[sup_broadcast]
redis = redis://localhost:6379
supervisor_url = http://localhost:19001
